# bridge/com.marginalia.bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- marginalia LaunchAgent template — do NOT edit this file directly.
     Run bridge/install.sh to generate and install the populated version. -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.marginalia.bridge</string>

    <key>ProgramArguments</key>
    <array>
        <string>{{PYTHON}}</string>
        <string>{{BRIDGE_SERVER}}</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <!-- Set your API key for whichever provider you use: -->
        <!-- <key>MARGINALIA_OPENAI_API_KEY</key><string>sk-...</string> -->
        <!-- <key>MARGINALIA_ANTHROPIC_API_KEY</key><string>sk-ant-...</string> -->
        <!-- <key>MARGINALIA_AWS_PROFILE</key><string>your-aws-profile</string> -->
        <key>MARGINALIA_MODEL_ID</key>
        <string>openai:gpt-4o</string>
        <key>MARGINALIA_VAULT</key>
        <string>{{VAULT_PATH}}</string>
        <key>MARGINALIA_PORT</key>
        <string>7731</string>
        <!-- Uncomment to require a shared secret in the KOReader plugin: -->
        <!-- <key>MARGINALIA_TOKEN</key><string>your-secret</string> -->
        <key>PATH</key>
        <string>{{PATH_DIRS}}</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>{{LOG_PATH}}</string>
    <key>StandardErrorPath</key>
    <string>{{LOG_PATH}}</string>

    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
